# deploy/supervisord.conf
; The full Tideo stack, co-located in one HF Space container. Every service runs as its own process
; (uid 1000), all over localhost. Only the API is exposed (port 7860). Boot order: stores -> brokers -> app.
; startretries is high so boot-order races (a client starting before its store accepts) self-heal via
; autorestart instead of going FATAL.
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/data/logs
environment=PATH="/opt/venv/bin:/opt/kafka/bin:/usr/sbin:/usr/bin:/bin",HOME="/data/home"

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---- stores ----
[program:redis]
command=/app/deploy/scripts/redis.sh
priority=10
autorestart=true
startsecs=3
startretries=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:postgres]
command=/app/deploy/scripts/postgres.sh
priority=10
autorestart=true
startsecs=6
startretries=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- brokers ----
[program:rabbitmq]
command=/app/deploy/scripts/rabbitmq.sh
priority=20
autorestart=true
startsecs=12
startretries=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:kafka]
command=/app/deploy/scripts/kafka.sh
priority=20
autorestart=true
startsecs=12
startretries=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- app ----
[program:api]
command=uvicorn app.api.main:app --host 0.0.0.0 --port 7860
directory=/app
priority=40
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:dispatcher]
command=python -m app.dispatcher
directory=/app
priority=50
autorestart=true
startsecs=5
startretries=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:audit]
command=python -m app.dispatcher.audit
directory=/app
priority=50
autorestart=true
startsecs=5
startretries=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker-heavy]
command=celery -q -A app.workers.celery_app worker -Q transcode --concurrency=2 --prefetch-multiplier=1 -n heavy@%%h --loglevel=info
directory=/app
priority=50
autorestart=true
startsecs=8
startretries=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker-fast]
command=celery -q -A app.workers.celery_app worker -Q inspect,package,transcribe,cleanup --concurrency=4 --prefetch-multiplier=2 -n fast@%%h --loglevel=info
directory=/app
priority=50
autorestart=true
startsecs=8
startretries=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:beat]
command=celery -q -A app.workers.celery_app beat --schedule=/data/celerybeat-schedule --loglevel=info
directory=/app
priority=60
autorestart=true
startsecs=5
startretries=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
